# etc/s6-overlay/s6-rc.d/adsbexchange-hassio-sensors/run
#!/command/with-contenv bash

if [ "$SERVICE_ENABLE_ADSBEXCHANGE" != "false" ]; then
    set -eo pipefail

    /usr/local/share/adsbexchange/hassio-sensors 2>&1 | mawk -W interactive '{printf "%c[35m[adsbexchange-hassio]%c[0m %s\n", 27, 27, $0}'
    # awk -W interactive ...  (prefix log messages with color and "[adsbexchange-hassio]")
else
    tail -f /dev/null
fi
