# etc/s6-overlay/s6-rc.d/adsbfi-hassio-sensors/run
#!/command/with-contenv bash

if [ "$SERVICE_ENABLE_ADSBFI" != "false" ]; then
    set -eo pipefail

    /usr/local/share/adsbfi/hassio-sensors 2>&1 | mawk -W interactive '{printf "%c[35m[adsbfi-hassio]%c[0m %s\n", 27, 27, $0}'
    # awk -W interactive ...  (prefix log messages with color and "[adsbfi-hassio]")
else
    tail -f /dev/null
fi
